# Library/LaunchAgents/com.thevetat.deskflow-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.thevetat.deskflow-server</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/Deskflow.app/Contents/MacOS/deskflow-core</string>
		<string>server</string>
		<string>--settings</string>
		<string>/Users/thevetat/Library/Deskflow/Deskflow.conf</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>KeepAlive</key>
	<false/>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
	<key>StandardOutPath</key>
	<string>/tmp/deskflow-server.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/deskflow-server.err</string>
</dict>
</plist>
